# deb/postinst
#!/bin/bash
systemctl enable tisbackup_huey
systemctl enable tisbackup_gui

if [ ! -f /usr/bin/tisbackup ]; then
    ln -s /opt/tisbackup/tisbackup.py  /usr/bin/tisbackup
fi
