# deb/postinst
#!/bin/bash
systemctl enable tisbackup_huey
systemctl enable tisbackup_gui

