# deb/postinst
#!/bin/bash
systemctl enable tisbackup_huey
systemctl enable tisbackup_gui
ln -s /opt/tisbackup/tisbackup.py  /usr/bin/tisbackup

